# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:osint-adapter]
command=python3 /app/osint-adapter/main.py
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/osint-adapter.log
stderr_logfile=/var/log/supervisor/osint-adapter-err.log
priority=5
environment=A0_PATH="/a0",API_KEY_OPENROUTER="%(ENV_API_KEY_OPENROUTER)s"

[program:spring-boot]
command=java -jar -Dspring.profiles.active=%(ENV_SPRING_PROFILES_ACTIVE)s -DDATA_DIR=/app/data /app/app.jar
autorestart=true
startsecs=10
stdout_logfile=/var/log/supervisor/spring-boot.log
stderr_logfile=/var/log/supervisor/spring-boot-err.log
priority=10
environment=JWT_SECRET="%(ENV_JWT_SECRET)s",DATA_DIR="/app/data",AUGUST_ADAPTER_URL="http://localhost:8000",API_KEY_OPENROUTER="%(ENV_API_KEY_OPENROUTER)s"
