# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:osint-adapter]
command=python3 /app/osint-adapter/main.py
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/osint-adapter.log
stderr_logfile=/var/log/supervisor/osint-adapter-err.log
priority=5

[program:spring-boot]
command=java -jar -Dspring.profiles.active=%(ENV_SPRING_PROFILES_ACTIVE)s -DDATA_DIR=/app/data /app/app.jar
autorestart=true
startsecs=10
stdout_logfile=/var/log/supervisor/spring-boot.log
stderr_logfile=/var/log/supervisor/spring-boot-err.log
priority=10
