# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:osint-adapter]
command=python3 /app/osint-adapter/main.py
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/osint-adapter.log
stderr_logfile=/var/log/supervisor/osint-adapter-err.log
priority=5
environment=
    A0_PATH="/a0",
    API_KEY_OPENROUTER="%(ENV_API_KEY_OPENROUTER)s",
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s",
    GOOGLE_CSE_ID="%(ENV_GOOGLE_CSE_ID)s",
    PROXYCURL_API_KEY="%(ENV_PROXYCURL_API_KEY)s",
    PDL_API_KEY="%(ENV_PDL_API_KEY)s",
    HUNTER_API_KEY="%(ENV_HUNTER_API_KEY)s",
    OSINT_ENGINE="%(ENV_OSINT_ENGINE)s"

[program:spring-boot]
command=java -jar
    -Dspring.profiles.active=%(ENV_SPRING_PROFILES_ACTIVE)s
    -DDATA_DIR=/app/data
    /app/app.jar
autorestart=true
startsecs=10
stdout_logfile=/var/log/supervisor/spring-boot.log
stderr_logfile=/var/log/supervisor/spring-boot-err.log
priority=10
environment=
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    DATA_DIR="/app/data",
    APP_BASE_URL="%(ENV_APP_BASE_URL)s",
    AUGUST_ADAPTER_URL="http://localhost:8000",
    OSINT_ENGINE="%(ENV_OSINT_ENGINE)s",
    AGENT_ZERO_URL="%(ENV_AGENT_ZERO_URL)s",
    API_KEY_OPENROUTER="%(ENV_API_KEY_OPENROUTER)s",
    MAIL_ENABLED="%(ENV_MAIL_ENABLED)s",
    MAIL_HOST="%(ENV_MAIL_HOST)s",
    MAIL_PORT="%(ENV_MAIL_PORT)s",
    MAIL_USERNAME="%(ENV_MAIL_USERNAME)s",
    MAIL_PASSWORD="%(ENV_MAIL_PASSWORD)s",
    MAIL_FROM="%(ENV_MAIL_FROM)s",
    STRIPE_SECRET_KEY="%(ENV_STRIPE_SECRET_KEY)s",
    STRIPE_WEBHOOK_SECRET="%(ENV_STRIPE_WEBHOOK_SECRET)s",
    STRIPE_PRICE_ID_STARTER="%(ENV_STRIPE_PRICE_ID_STARTER)s",
    STRIPE_PRICE_ID_PROFESSIONAL="%(ENV_STRIPE_PRICE_ID_PROFESSIONAL)s"
