# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:spring-boot]
command=java -jar -Dspring.profiles.active=prod -DDATA_DIR=/data /app/app.jar
autorestart=true
stdout_logfile=/var/log/supervisor/spring-boot.log
stderr_logfile=/var/log/supervisor/spring-boot-err.log
priority=10

[program:osint-adapter]
command=python3 /app/osint-adapter/main.py
autorestart=true
stdout_logfile=/var/log/supervisor/osint-adapter.log
stderr_logfile=/var/log/supervisor/osint-adapter-err.log
priority=5
environment=A0_PATH="/a0"
